# deadline-sync/com.changliu.deadline-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.changliu.deadline-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/changliu/.conda/envs/deadline-sync/bin/python</string>
        <string>/Users/changliu/my-todo/deadline-sync/main.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/changliu/my-todo/deadline-sync</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/changliu/Library/Logs/deadline-sync.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/changliu/Library/Logs/deadline-sync-error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
